# init/supervisord.conf
[supervisor]
ndaemon = True

[proagram:rasa_default_server]
command: rasa run

[program:rasa_action_server]
command: rasa run actions
